# recipes/chumby-udev-rules/files/50-firmware-compat.rules
SUBSYSTEM=="compat_firmware", ACTION=="add", RUN+="firmware --firmware=$env{FIRMWARE} --devpath=$env{DEVPATH}"
